# crontab
0 0 * * * /home/ubuntu/import-officialevent-bat/scripts/daily_batch.sh
