# systemd/import-officialevent-bat.timer
[Unit]
Description=import-officialevent-bat

[Timer]
OnCalendar=*-*-* 0:15:00
Persistent=true

[Install]
WantedBy=timers.target